# etc/cron.daily/1-paludis-sync
#!/bin/bash

pushd /tmp
sudo -u paludisbuild -g paludisbuild cave sync
sudo -u paludisbuild -g paludisbuild cave manage-search-index
popd
